# configs/vpnserver
#!/bin/sh
# chkconfig: 2345 99 01
# description: SoftEther VPN Server
DAEMON=/usr/local/vpnserver/vpnserver
LOCK=/var/lock/subsys/vpnserver

IPTABLESBIN=/usr/sbin/iptables
IFCONFIG=/sbin/ifconfig
IPBIN=/sbin/ip

SERVER_IP=124.124.124.124

SERVER_NIC=$(ip -4 route ls | grep default | grep -Po '(?<=dev )(\S+)' | head -1)

TAP_INTERFACE=tap_in
TAP_ADDR=9.9.9.1
TAP_NETWORK=9.9.9.0/24
TAP_STATIC=11.11.111.50/24

test -x $DAEMON || exit 0
case "$1" in

start)
echo " Setting IP Tables"
$IPTABLESBIN -t nat -D POSTROUTING -s $TAP_NETWORK -o $SERVER_NIC -j MASQUERADE
$IPTABLESBIN -t nat -D POSTROUTING -s $TAP_NETWORK -o $SERVER_NIC -j MASQUERADE
$IPTABLESBIN -t nat -D POSTROUTING -s $TAP_NETWORK -o $SERVER_NIC -j MASQUERADE
sleep 1
$IPTABLESBIN -t nat -A POSTROUTING -s $TAP_NETWORK -o $SERVER_NIC -j MASQUERADE
$DAEMON start
touch $LOCK
sleep 1

# FOR EUROPE SERVER
#$IFCONFIG $TAP_INTERFACE $TAP_ADDR 

# FOR DOMESTIC [IRAN] SERVER
$IPBIN addr add  $TAPSTATIC brd + dev $TAP_INTERFACE

/etc/init.d/dnsmasq restart
#$IPTABLESBIN -t nat -A POSTROUTING -s $TAP_NETWORK -o $SERVER_NIC -j MASQUERADE
;;

stop)
$IPTABLESBIN -t nat -D POSTROUTING -s $TAP_NETWORK -o $SERVER_NIC -j MASQUERADE
$IPTABLESBIN -t nat -D POSTROUTING -s $TAP_NETWORK -o $SERVER_NIC -j MASQUERADE
$IPTABLESBIN -t nat -D POSTROUTING -s $TAP_NETWORK -o $SERVER_NIC -j MASQUERADE
$DAEMON stop
rm $LOCK
;;

restart)
$DAEMON stop
$IPTABLESBIN -t nat -D POSTROUTING -s $TAP_NETWORK -o $SERVER_NIC -j MASQUERADE
$IPTABLESBIN -t nat -D POSTROUTING -s $TAP_NETWORK -o $SERVER_NIC -j MASQUERADE
$IPTABLESBIN -t nat -D POSTROUTING -s $TAP_NETWORK -o $SERVER_NIC -j MASQUERADE
sleep 1
$DAEMON start
$IPTABLESBIN -t nat -A POSTROUTING -s $TAP_NETWORK -o $SERVER_NIC -j MASQUERADE
sleep 1

/sbin/ip addr add  $TAPSTATIC brd + dev $TAP_INTERFACE
#$IFCONFIG $TAP_INTERFACE $TAP_ADDR
;;

*)
echo "Usage: $0 {start|stop|restart}"
exit 1
esac
exit 0
